# ansible/scripts/concourse/concourse-worker.service
[Unit]
Description=Concourse CI worker process
After=concourse-web.service

[Service]
User=root
Restart=on-failure
ExecStart=/usr/local/bin/concourse-worker

[Install]
WantedBy=multi-user.target
